# ts/build/packages/ica_usb/build/extra/etc/systemd/system/ctxusbd.service
[Unit]
Description=CTX USB Daemon
After=ica.init.service ctxcwalogd.service

[Service]
Type=forking
ExecStart=/etc/init.d/ctxusbd start
ExecStop=/etc/init.d/ctxusbd stop
Restart=on-failure

[Install]
WantedBy=multi-user.target
