# ts/build/packages/persistent-files/etc/systemd/system/persistent-files.service
[Unit]
Description=Monitor/Backup/Restore system critical files
DefaultDependencies=no
After=var-prstnt.mount fastboot.service
Before=display-manager.service tsinit.target
ConditionPathIsMountPoint=/var/prstnt

[Service]
Type=simple
ExecStartPre=/etc/init.d/persistent-files restore
ExecStart=/etc/init.d/persistent-files monitor
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
