# ts/build/packages/networkmanager/build/extra/etc/systemd/system/NetworkManager.service.d/overrides.conf
[Service]
ExecStartPre=/etc/init.d/nmcreatesysconn
