# ts/build/packages/systemd/build/extra/etc/systemd/system/systemd-tmpfiles-reload.service
[Unit]
Description=Restart systemd-tmpfiles after fastboot
After=fastboot.service
Requires=fastboot.service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart systemd-tmpfiles-setup.service

[Install]
WantedBy=multi-user.target
