# ts/build/packages/systemd/build/extra/etc/systemd/system/systemd-vconsole-setup-reload.service
[Unit]
Description=Restart Virtual Console Setup Service
After=tsinit.target profile-setup.service net.profile-setup.service

[Service]
Type=oneshot
ExecStart=systemctl restart systemd-vconsole-setup.service

[Install]
WantedBy=tsinit.target
