# ts/build/packages/autonet/etc/systemd/system/choose-network.service
[Unit]
Description=Choose from available network connections
After=post-udev.service	network-pre.target dbus.service autonet.service
Before=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/choose-network
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
