# ts/build/packages/horizon/build/extra/etc/systemd/system/horizon-USBArbitrator.service.d/overrides.conf
[Unit]
After=pkg.service tsinit.target

[Install]
WantedBy=
WantedBy=tsinit.target
