# ts/build/packages/git/build/extra/etc/systemd/system/git-init.service
[Unit]
Description=ThinStation Git Init
After=profile-setup.service tsinit.target
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/git-init
SyslogIdentifier=thinstation

[Install]
WantedBy=tsinit.target
