# ts/build/packages/networkmanager/build/extra/etc/systemd/system/hostname-reconfigure.service
[Unit]
Description=Run nmcreatesyscon script
After=profile-setup.service net.profile-setup.service
After=tsinit.target
Before=display-manager.service lightdm.service

[Service]
Type=oneshot
ExecStart=/etc/init.d/nmcreatesysconn reconfigure
RemainAfterExit=true

[Install]
WantedBy=tsinit.target
