# ts/build/packages/sssd/build/extra/etc/systemd/system/oddjobd.service.d/overrides.conf
[Unit]
After=tsinit.target sssd.service
Requires=sssd.service

[Install]
WantedBy=
WantedBy=tsinit.target

