# ts/build/packages/base/etc/systemd/system/crond.service
[Unit]
Description=Command Scheduler
After=tsinit.target crontab.service

[Service]
ExecStart=/sbin/crond -f
KillMode=process

[Install]
WantedBy=tsinit.target
