# ts/build/packages/base/etc/systemd/system/crontab.service
[Unit]
Description=ThinStation Crontab Initialize
Before=crond.service
After=tsinit.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/crontab
SyslogIdentifier=thinstation

[Install]
WantedBy=tsinit.target
