# ts/build/packages/base/etc/systemd/system/fastboot.service
[Unit]
Description=Mount Fastboot Data
After=network-online.target
Before=tsinit.target display-manager.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/fastboot
SyslogIdentifier=fastboot
TimeoutStartSec=180s

[Install]
WantedBy=multi-user.target
