# ts/build/packages/base/etc/systemd/system/persistent-dirs.service
[Unit]
Description=Monitor/Backup/Restore system critical dirs
DefaultDependencies=no
After=var-prstnt.mount fastboot.service
Before=display-manager.service tsinit.target
ConditionPathIsMountPoint=/var/prstnt

[Service]
Type=simple
ExecStartPre=/etc/init.d/persistent-dirs restore
ExecStart=/etc/init.d/persistent-dirs monitor
Restart=failure
RestartSec=10

[Install]
WantedBy=multi-user.target
