# ts/build/packages/base/etc/systemd/system/pkg.service
[Unit]
Description=ThinStation Pkg files
ConditionPathIsReadWrite=/etc
After=tsinit.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/pkg init
SyslogIdentifier=thinstation

[Install]
WantedBy=tsinit.target

