# ts/build/packages/base/etc/systemd/system/profile-setup.service
[Unit]
Description=ThinStation Profile Setup
ConditionPathIsReadWrite=/etc
After=persistent-files.service persistent-dirs.service fastboot.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/profile-setup
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
