# ts/build/packages/base/etc/systemd/system/profile-setup.service
[Unit]
Description=ThinStation Profile Setup
ConditionPathIsReadWrite=/etc
After=tsinit.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/profile-setup
SyslogIdentifier=thinstation

[Install]
WantedBy=tsinit.target
