# ts/build/packages/base/etc/systemd/system/session-setup.service
[Unit]
Description=ThinStation Session Configure
After=tsinit.target session.service
Before=display-manager.service lightdm.service
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/session-setup
SyslogIdentifier=thinstation

[Install]
WantedBy=tsinit.target
