# ts/build/packages/base/etc/systemd/system/tsinit.target
[Unit]
Description=ThinStation Initialization Target

After=fastboot.service
Before=display-manager.service lightdm.service

[Install]
WantedBy=multi-user.target
