# ts/build/packages/base/etc/systemd/system/vmtouch.service
[Unit]
Description=Portable file system cache diagnostics and control
After=tsinit.target

[Service]
Type=oneshot
ExecStart=/bin/vmtouch -ft /lib64
RemainAfterExit=yes

[Install]
WantedBy=tsinit.target
