# ts/build/packages/spice-vdagent/build/extra/etc/systemd/system/9p-qemu-mounts.service
[Unit]
Description=Mount shared folder between qemu host and guest
ConditionVirtualization=qemu
After=post-udev.service fastboot.service

[Service]
Type=oneshot
ExecStart=/etc/init.d/9p-mounts
RemainAfterExit=yes


[Install]
WantedBy=multi-user.target
