# ts/build/packages/netfiles/etc/systemd/system/network-files.service
[Unit]
Description=ThinStation Network Files
After=profile-setup.service network-online.target tsinit.target
Before=net.profile-setup.service avahi-daemon.service
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/network-files
SyslogIdentifier=thinstation

[Install]
WantedBy=tsinit.target
