# ts/build/packages/vmview/build/extra/etc/systemd/system/ftscanhv.service.d/overrides.conf
[Unit]
After=pkg.service

[Install]
WantedBy=tsinit.target
