# ts/build/packages/xorg7vnc/etc/systemd/system/x11vnc.service
[Unit]
Description=ThinStation x11vnc systemd unit
After=profile-setup.service pkg.service tsinit.target
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/x11vnc init
SyslogIdentifier=thinstation

[Install]
WantedBy=tsinit.target
