# ts/build/packages/firefox/build/extra/etc/systemd/system/firefox-init.service
[Unit]
Description=Configure Firefox Settings
After=profile-setup.service net.profile-setup.service pkg.service tsinit.target
Before=session-setup.service
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/firefox-init
SyslogIdentifier=thinstation

[Install]
WantedBy=tsinit.target
