# ts/build/packages/ca-certificates/build/extra/etc/systemd/system/update-ca-trust.path
[Unit]
Description=Watch /etc/pki/ca-trust/source/anchors for changes

[Path]
PathExistsGlob=/etc/pki/ca-trust/source/anchors/*
Unit=update-ca-trust.service

[Install]
WantedBy=multi-user.target
