# ts/build/packages/ca-certificates/build/extra/etc/systemd/system/update-ca-trust.service
[Unit]
Description=Update CA Trust if files are present in anchors

[Service]
Type=oneshot
ExecStart=/bin/update-ca-trust
