# ts/build/packages/openssh/build/extra/etc/systemd/system/sshd.service.d/overrides.conf
[Unit]
After=tsinit.target networkmanager-reload.service

[Install]
WantedBy=
WantedBy=tsinit.target
