# ts/build/packages/horizon/build/extra/etc/systemd/system/horizon-init.service
[Unit]
Description=Configure Omnissa Horizon Client Settings
DefaultDependencies=no
Conflicts=shutdown.target
After=profile-setup.service net.profile-setup.service pkg.service session.service tsinit.target
Before=session-setup.service horizon-USBArbitrator.service
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/horizon-init
SyslogIdentifier=thinstation

[Install]
WantedBy=tsinit.target
