# ts/build/packages/open-vm-tools/build/extra/etc/systemd/system/vmware-mounts.service
[Unit]
Description=Mount VMware Shared Folders Dynamically
After=vmtoolsd.service
ConditionVirtualization=vmware

[Service]
Type=oneshot
ExecStart=/etc/init.d/vmware-mounts start
ExecStop=/etc/init.d/vmware-mounts stop
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
