# ts/build/packages/nftables/build/extra/etc/systemd/system/nftables.service.d/overrides.conf
[Service]
ExecStartPost=/etc/init.d/firewall
