# ts/build/packages/devstation/etc/systemd/system/networkmanager-reload.service
[Unit]
Description=Restart NetworkManager After Clearing State
After=persistent-files.service persistent-dirs.service profile-setup.service
Requires=persistent-files.service persistent-dirs.service
ConditionPathExists=/var/prstnt/etc/NetworkManager

[Service]
Type=oneshot
ExecStartPre=/bin/rm -rf /run/NetworkManager
#ExecStart=/bin/nmcli connection reload
ExecStart=systemctl restart NetworkManager
ExecStartPost=nm-online -s -q

[Install]
WantedBy=multi-user.target
