# ts/build/packages/devstation/etc/systemd/system/dnsmasq.service.d/overrides.conf
[Unit]
After=persistent-files.service restart-networkmanager.service NetworkManager.service
ConditionPathExists=/etc/nftables/nftables.d/dnsmasq.conf

[Service]
ExecStartPost=systemctl restart nftables

[Install]
WantedBy=tsinit.target

