# ts/build/packages/devstation/etc/systemd/system/dnsmasq.service.d/overrides.conf
[Unit]
After=persistent-files.service networkmanager-reload.service NetworkManager.service firewalld.service
ConditionPathExists=/etc/dnsmasq.conf

[Install]
WantedBy=tsinit.target

