# ts/build/packages/devstation/etc/systemd/system/dnsmasq.service.d/overrides.conf
[Unit]
After=tsinit.target networkmanager-reload.service
ConditionPathExists=/etc/dnsmasq.conf

[Install]
WantedBy=
WantedBy=tsinit.target

